# docker/python/etc/supervisor/conf.d/uwsgi.conf
[program:uwsgi]
directory=/var/www/app
command=/usr/local/bin/uwsgi --ini /var/www/app/uwsgi.ini
numprocs=1
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/uwsgi.log
